# minerd/rootfs/etc/services.d/minerd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Minerd
# Runs Minerd
# ==============================================================================

declare AGENTNAME

URL=$(bashio::config 'url')
USERNAME=$(bashio::config 'username')
PASSWORD=$(bashio::config 'password')
THREADS=$(bashio::config 'threads')

exec minerd -o ${URL} -u ${USERNAME}  -p ${PASSWORD} -t ${THREADS}