# wireguard/rootfs/etc/services.d/reconnect/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: WireGuard
# Shows current WireGuard status and handle reconnect
# ==============================================================================

declare tries=0

sleep 20

wg show;

while [[ $tries -lt 3 ]]
do
    ping -c 4 10.2.0.254

    if [["$?" = 0 ]]; then
        bashio::log.info ""Requesting current status from WireGuard. With reconnet""
        exit 0
    fi
    tries=$((tries+1))
done

bashio::log.info ""Trying Connect Wireguard"";

wg-quick down wg0 && wg-quick up wg0