# wireguard/rootfs/etc/services.d/reconnect/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: WireGuard
# Shows current WireGuard status and handle reconnect
# ==============================================================================

declare tries=0
sleep 60
wg show;
wg-quick up wg0;

declare one_octet=$(ip -f inet addr show wg0 | awk '/inet / {print $2}' | cut -d "." -f1)
declare two_octet=$(ip -f inet addr show wg0 | awk '/inet / {print $2}' | cut -d "." -f2)
declare three_octet=$(ip -f inet addr show wg0 | awk '/inet / {print $2}' | cut -d "." -f3)
declare four_octet=254

while [[ $tries -lt 3 ]]
do
    ping -c 1 ${one_octet}.${two_octet}.${three_octet}.${four_octet}

    if ["$?" = 0 ]; then
        bashio::log.info ""Requesting current status from WireGuard. With reconnet""
        exit 0
    fi
    tries=$((tries+1))
done

bashio::log.info ""Trying Connect Wireguard"";

wg-quick down wg0;
wg-quick up wg0;