# connect/rootfs/etc/services.d/connect/finish
#!/command/with-contenv bashio

# ==============================================================================
# Unicontrol Add-on: Connect
# ==============================================================================

declare interface_ifname=$(bashio::config "interface.ifname")

wg-quick down ${interface_ifname}